# system-configs/crontab
*/10 * * * * /home/relive/scripts/get-releases.sh
0 * * * * /home/relive/scripts/get-fahrplan.sh
*/2 * * * * perl /home/relive/hls-relive/genpage.pl
